# roles/plexmediaserver/files/common/plexmediaserver.service.d/restrict.conf
[Service]
ReadOnlyPaths=/
ReadWritePaths=/var/lib/plex /tmp /media/movies